# scripts/model/crontab
# Legacy cron alternative. The VPS systemd timer is preferred because it limits CPU usage.
# This job downloads data and evaluates a candidate; daily mode never promotes it.
15 3 * * * SUPER_POKER_PYTHON=/root/super_poker_3/.venv/bin/python SUPER_POKER_DATA_DIR=/root/super_poker_3/data/raw /root/super_poker_3/scripts/model/vps_daily_candidate.sh >> /root/super_poker_3/logs/daily-candidate-cron.log 2>&1

# Do not schedule weekly_learning.sh or cycle_learning.sh unattended on a competition miner.
